# etc/cron.d/0100-samvirke
0 1 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:samvirke-latest reload/db-datadir:samvirke-latest) 2>&1 | /usr/bin/logger -t datadir-preinit:samvirke
