# etc/cron.d/0200-spejderdk
0 2 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:spejderdk-latest reload/db-datadir:spejderdk-latest ../jobs/spejderdk/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:spejderdk
