# etc/cron.d/0230-spejderneslejr
30 2 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/spejderneslejr:db-data-test reload/spejderneslejr:db-datadir-test ../jobs/spejderneslejr/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:spejderneslejr