# etc/cron.d/0300-helbredsprofilen
0 3 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh helbredsprofilen-latest ../jobs/helbredsprofilen/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:helbredsprofilen
