# etc/cron.d/0330-msf
30 3 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:msf-latest reload/db-datadir:msf-latest ../jobs/msf/datadir-preinit-sanitize.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:msf-latest reload/db-datadir:msf-latest ../jobs/msf/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:msf
