# etc/cron.d/0400-spejderneslejr
30 4 * * * root (. /etc/container_environment.sh \
 && cd /opt/datadir-preinit \
 && BASE_IMAGE=mariadb:10 ./preinit-tag.sh reload/spejderneslejr:db-data-prod reload/spejderneslejr:mariadb-10-prod ../jobs/spejderneslejr/datadir-preinit-sanitize.sql \
 && BASE_IMAGE=mariadb:10 ./preinit-tag.sh reload/spejderneslejr:db-data-test reload/spejderneslejr:mariadb-10-test ../jobs/spejderneslejr/datadir-preinit-sanitize.sql \
 && BASE_IMAGE=mariadb:10 ./preinit-tag.sh reload/spejderneslejr:db-data-stage reload/spejderneslejr:mariadb-10-stage ../jobs/spejderneslejr/datadir-preinit-sanitize.sql \
 ) 2>&1 | /usr/bin/logger -t datadir-preinit:spejderneslejr