# etc/cron.d/0415-dj
# Dump is made 4:00
15 4 * * * root time (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:dj-latest reload/db-datadir:dj-latest /opt/jobs/general-d8/reset-admin-pass.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:dj-latest reload/db-datadir:dj-latest /opt/jobs/general-d8/reset-admin-pass.sql && ./preinit-tag.sh reload/db-data:dj-ajks-latest reload/db-datadir:dj-ajks-latest /opt/jobs/general-d8/reset-admin-pass.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:dj-ajks-latest reload/db-datadir:dj-ajks-latest /opt/jobs/general-d8/reset-admin-pass.sql ) 2>&1 | /usr/bin/logger -t datadir-preinit:dj
