# etc/cron.d/0500-ddsdk2016
# Dump is made 4:30
00 5 * * * root time (. /etc/container_environment.sh && cd /opt/datadir-preinit && (./preinit-tag.sh reload/db-data:ddsdk2016-master-latest reload/db-datadir:ddsdk2016-master-latest /opt/jobs/general-d8/reset-admin-pass.sql; FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:ddsdk2016-master-latest reload/db-datadir:ddsdk2016-master-latest /opt/jobs/general-d8/reset-admin-pass.sql; ./preinit-tag.sh reload/db-data:ddsdk2016-develop-latest reload/db-datadir:ddsdk2016-develop-latest /opt/jobs/general-d8/reset-admin-pass.sql; FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:ddsdk2016-develop-latest reload/db-datadir:ddsdk2016-develop-latest /opt/jobs/general-d8/reset-admin-pass.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:ddsdk2016
