# etc/cron.d/0500-ddsdk2016
0 5 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:ddsdk2016-latest reload/db-datadir:ddsdk2016-latest ../jobs/general-d8/reset-admin-pass.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:ddsdk2016
