# etc/cron.d/0500-ida
0 5 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && NO_VOLUME=1 ./preinit-tag.sh reload/db-data:ida-latest reload/db-datadir:ida-latest ../jobs/ida/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:ida
