# etc/cron.d/0530-ida
# Dump is made 5:00
# Notice, this dump is big so allow at least '30 mins for it to complete.
30 5 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:ida-latest reload/db-datadir:ida-latest ../jobs/ida/datadir-preinit-sanitize.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:ida-slim-latest reload/db-datadir:ida-slim-latest ../jobs/ida/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:ida