# etc/cron.d/0615-spejderdk
# Dump is made 5:30
15 6 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:spejderdk-master-latest reload/db-datadir:spejderdk-master-latest ../jobs/spejderdk/datadir-preinit-sanitize.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:spejderdk-master-latest reload/db-datadir:spejderdk-master-latest ../jobs/spejderdk/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:spejderdk
25 6 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && ./preinit-tag.sh reload/db-data:spejderdk-develop-latest reload/db-datadir:spejderdk-develop-latest ../jobs/spejderdk/datadir-preinit-sanitize.sql && FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:spejderdk-develop latest reload/db-datadir:spejderdk-develop-latest ../jobs/spejderdk/datadir-preinit-sanitize.sql) 2>&1 | /usr/bin/logger -t datadir-preinit:spejderdk
