# etc/cron.d/0700-ddsdk2016
0 7 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && m) 2>&1 | /usr/bin/logger -t datadir-preinit:ddsdk2016
