# etc/cron.d/0800-lokalbolig
# Dump is made 7:30
00 8 * * * root (. /etc/container_environment.sh && cd /opt/datadir-preinit && (./preinit-tag.sh reload/db-data:lokalbolig-master-latest reload/db-datadir:lokalbolig-master-latest /opt/jobs/general-d8/reset-admin-pass.sql; FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:lokalbolig-master-latest reload/db-datadir:lokalbolig-master-latest /opt/jobs/general-d8/reset-admin-pass.sql; ./preinit-tag.sh reload/db-data:lokalbolig-develop-latest reload/db-datadir:lokalbolig-develop-latest /opt/jobs/general-d8/reset-admin-pass.sql; FULL_DB_IMAGE=1 ./preinit-tag.sh reload/db-data:lokalbolig-develop-latest reload/db-datadir:lokalbolig-develop-latest /opt/jobs/general-d8/reset-admin-pass.sql)) 2>&1 | /usr/bin/logger -t datadir-preinit:lokalbolig
